# supervisord.php.conf
[program:php5-fpm]
command=php5-fpm --nodaemonize
priority=30
stderr_logfile=/dev/null
stdout_logfile=/dev/null
