# supervisord.php.conf
[program:php-fpm]
autorestart = true
command = /usr/sbin/php-fpm8.1 --nodaemonize
priority = 30
stderr_logfile = /dev/null
stdout_logfile = /dev/null
