# supervisord.nginx.conf
[program:nginx]
command = nginx
priority = 40
stderr_logfile = /dev/null
stdout_logfile = /dev/null
