# supervisord.php.conf
[program:php-fpm]
command = php-fpm7.4 --nodaemonize
priority = 30
stderr_logfile = /dev/null
stdout_logfile = /dev/null
