# supervisord.php.conf
[program:php7.0-fpm]
command=php-fpm7.0 --nodaemonize
priority=30
stderr_logfile=/dev/null
stdout_logfile=/dev/null
